# supervisord.conf
[program:refresh]
command=python refresh.py

[program:serve]
command=flask run --with-threads

[program:screenshot]
command=node screenshot.js

[supervisord]
logfile = supervisord.log

[unix_http_server]
file=supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://supervisord.sock
